# etc/service/sabnzbd/run
#! /bin/sh

# start sabnzbd
# running in the foreground ill cause the container to stop when sabnzbd stops
/usr/bin/sabnzbdplus --config-file /config --server :8080


